# src/private_dot_config/systemd/user/tmux-htop-session.service
[Unit]
Description=Start tmux session with htop
After=default.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tmux-htop-session
RemainAfterExit=yes

[Install]
WantedBy=default.target

